# template_compiles/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl
#------------------------------------------------------------------------
# Compiled template generated by the Template Toolkit version 3.100
#------------------------------------------------------------------------

Template::Document->new({
    METADATA => {
        'modtime' => '1728930269',
        'name' => '/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl',
    },
    BLOCK => sub {
    my $context = shift || die "template sub called without context\n";
    my $stash   = $context->stash;
    my $output  = '';
    my $_tt_error;

    eval { BLOCK: {
$output .=  '#!/bin/bash
#
# ';
#line 3 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('service');
$output .=  '
#
# chkconfig: 2345 ';
#line 5 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('start_prio');
$output .=  ' ';
#line 5 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('stop_prio');
$output .=  '
# description: ';
#line 6 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('service');
$output .=  '
';
#line 15 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
if ($stash->get(['pidfile', 0, 'defined', 0]) && $stash->get(['pidfile', 0, 'length', 0]) > 0) {
$output .=  '
. /etc/init.d/cpfunctions

USER=';
#line 11 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('user');
$output .=  '
SERVICE=';
#line 12 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('service');
$output .=  '
PIDFILE=';
#line 13 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('pidfile');
$output .=  '
REGEX=\'';
#line 14 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('regex');
$output .=  '\'
';
}

$output .=  '
start() {';
#line 25 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
if ($stash->get(['pidfile', 0, 'defined', 0]) && $stash->get(['pidfile', 0, 'length', 0]) > 0) {
$output .=  "\n    already_running \$USER \$SERVICE \$PIDFILE \$REGEX; RUNNING=\$?\n    if [ 0 = \$RUNNING ]; then\n        ";
#line 21 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('start');
$output .=  '
    fi
';
}
else {
$output .=  '
    ';
#line 24 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('start');
$output .=  '
';
}

$output .=  '}

stop() {
    ';
#line 29 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('stop');
$output .=  '
}

status() {
    ';
#line 33 "/usr/local/cpanel/etc/init/scripts/templates/generic.tmpl"
$output .=  $stash->get('status');
$output .=  "\n}\n\ncase \"\$1\" in\n    start)\n        start\n        ;;\n    stop)\n        stop\n        ;;\n    restart)\n        stop\n        start\n        ;;\n    status)\n        status\n        ;;\nesac\n\n";
    } };
    if ($@) {
        $_tt_error = $context->catch($@, \$output);
        die $_tt_error unless $_tt_error->type eq 'return';
    }

    return $output;
},
    DEFBLOCKS => {

    },
});
